# 8.0/pgsql-nginx/master.ini
[supervisord]
logfile=/dev/fd/1
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=www-data
minfds=1024
minprocs=200

[program:crond]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command=crond -f -L 15
user=www-data
autostart=true

[program:php]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command = docker-php-entrypoint php-fpm
user=www-data

[program:nginx]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command=nginx -g 'pid /tmp/nginx.pid; daemon off;'
autostart=true
autorestart=true
user=www-data
priority=10
stdout_events_enabled=true
stderr_events_enabled=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[include]
files = /etc/supervisor.d/*.ini

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
