# base/nginx/supervisor.d/01_nginx.ini
[program:nginx]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
command=nginx -g 'pid /tmp/nginx.pid; daemon off; error_log /dev/stderr;'
autostart=true
autorestart=true
priority=10
