# base/nginx/supervisor.d/03_exporter_merger.ini
[program:exporter-merger]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on = php-fpm-prometheus-exporter,nginx-prometheus-exporter
environment=START_SCRIPT_CHECK_VARIABLE="ENABLE_PROMETHEUS_EXPORTER_RUNNER",START_SCRIPT_EXPECTED_VALUE="1"
command=/usr/local/bin/start-if-env-variable-is-set.sh
    /usr/bin/exporter-merger
    --listen-port 9090
    --url http://localhost:9190/metrics,http://localhost:9191/metrics
user=www-data
autostart=true
autorestart=true