# base/nginx/master.ini
[supervisord]
logfile=/dev/stderr
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=www-data
minfds=1024
minprocs=200

[program:crond]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
command=crond -f
user=www-data
autostart=true
autorestart=true

[program:php]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
command=docker-php-entrypoint php-fpm
user=www-data
autostart=true
autorestart=true
priority=10

[program:nginx]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
command=nginx -g 'pid /tmp/nginx.pid; daemon off; error_log /dev/stderr;'
autostart=true
autorestart=true
user=www-data
priority=10

[program:nginx-prometheus-exporter]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on = nginx
environment=START_SCRIPT_CHECK_VARIABLE="ENABLE_PROMETHEUS_EXPORTER_RUNNER",START_SCRIPT_EXPECTED_VALUE="1"
command=/usr/local/bin/start-if-env-variable-is-set.sh
    /usr/bin/nginx-prometheus-exporter
    -nginx.scrape-uri=http://127.0.0.1:8081/stub_status
    -web.listen-address=:9190
user=www-data
autostart=true
autorestart=true

[program:php-fpm-prometheus-exporter]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on = php
environment=START_SCRIPT_CHECK_VARIABLE="ENABLE_PROMETHEUS_EXPORTER_RUNNER",START_SCRIPT_EXPECTED_VALUE="1"
command=/usr/local/bin/start-if-env-variable-is-set.sh
    /usr/bin/php-fpm-exporter server
    --phpfpm.scrape-uri=tcp://127.0.0.1:9000/orchestration/php-fpm-status
    --phpfpm.fix-process-count=true
    --web.listen-address=:9191
user=www-data
autostart=true
autorestart=true

[program:exporter-merger]
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on = php-fpm-prometheus-exporter,nginx-prometheus-exporter
environment=START_SCRIPT_CHECK_VARIABLE="ENABLE_PROMETHEUS_EXPORTER_RUNNER",START_SCRIPT_EXPECTED_VALUE="1"
command=/usr/local/bin/start-if-env-variable-is-set.sh
    /usr/bin/exporter-merger
    --listen-port 9090
    --url http://localhost:9190/metrics,http://localhost:9191/metrics
user=www-data
autostart=true
autorestart=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9999

[include]
files = /etc/supervisor.d/*.ini

[inet_http_server]
port=127.0.0.1:9999
